# infrastructure/ansible/roles/job-runner/templates/fim-updateCurrentMatch.timer.j2
[Timer]
OnUnitInactiveSec={{ intervalSec }}s
OnBootSec={{ intervalSec }}s
Unit=fim-updateCurrentMatch.service
