# infrastructure/ansible/roles/fimadminapi-sync/templates/fim-adminapisync.timer.j2
[Timer]
OnUnitInactiveSec={{ intervalSec }}s
OnBootSec={{ intervalSec }}s
Unit=fim-adminapisync.service

[Install]
WantedBy=multi-user.target